# xivo-config/etc/supervisor/conf.d/xivo-agentd.conf
[program:xivo-agentd]
command = /usr/bin/xivo-agentd
redirect_stderr = true
stdout_logfile = /var/log/xivo-agentd.log
