# etc/systemd/system/rabbitmq-server.service.d/wazo.conf
[Service]
Environment=RABBITMQ_NODENAME=rabbit@localhost