# etc/systemd/system/rabbitmq-server.service.d/wazo.conf
[Service]
Environment=NODENAME=rabbit@localhost